# htpc/plex/crontab
# m h dom mon dow command
0 3 * * * /scripts/borg-backup.sh